# com.chibisafe.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chibisafe.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lajlev/Sites/tools/chibisafe-uploader/chibisafe_watcher.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/lajlev/Sites/tools/chibisafe-uploader/chibisafe_watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lajlev/Sites/tools/chibisafe-uploader/chibisafe_watcher.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/lajlev/Sites/tools/chibisafe-uploader</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
